# root/etc/services.d/epg-docker/run
#!/usr/bin/with-contenv bash

cd /app/epg || exit

exec \
  s6-setuidgid abc npm run serve
